# ec2/forge-worker/forge-secrets.service
[Unit]
Description=Forge worker secrets fetch (AWS Secrets Manager → /run/forge-worker.env)
# Network must be up; we're calling AWS API.
After=network-online.target
Wants=network-online.target
# Tag so consumer services can `Requires=forge-secrets.service`.

[Service]
Type=oneshot
# Runs as root because:
#   1. AWS SDK reads instance-metadata IAM role (works for any user).
#   2. We need to write /run/forge-worker.env (0600 root) and
#      /run/forge-worker/temporal-tls.key (0400 root).
#   3. Consumer services then read the file as their own user — systemd
#      handles that. Permissions are tight enough to keep secrets off
#      shell histories of non-root logged-in operators.
User=root
Group=root

# Default mode is strict — see forge-secrets-fetch.sh for details. Override
# to `soft` only during cutover so a bad IAM/SM state doesn't down a host
# that already has a working /run/forge-worker.env.
Environment=FORGE_SECRETS_MODE=strict
Environment=FORGE_SECRETS_ENV=prod
Environment=AWS_REGION=us-west-2

ExecStart=/opt/nexus/core/ec2/forge-worker/forge-secrets-fetch.sh

# Hardening — limits the blast radius if the script or aws CLI is exploited.
# /run is the only path the script writes to; no other parts of the FS need
# to be writable.
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
ReadWritePaths=/run
NoNewPrivileges=true
RestrictSUIDSGID=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX

# RemainAfterExit=yes means systemd treats this oneshot as "active" after
# successful run, which is what consumer Requires= expects to see.
RemainAfterExit=yes

StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
